# Utilities/input configuration/.config/autostart/wacom.desktop
[Desktop Entry]
Name=wacom configuration
Type=Application
Terminal=false
Exec=/usr/local/bin/wacom-second-button.sh
